# centreon-engine-1.4.15/scripts/deb/postinst.in
#!/bin/sh

# Abort if any command returns an error value
set -e

case "$1" in
    configure)
	if [ "@GROUP@" != "root" ]; then
	    if ! $(grep -E "^@GROUP@:" /etc/group > /dev/null ) ; then
	        groupadd -f "@GROUP@"
            fi
	fi

	if [ "@USER@" != "root" ]; then
            if ! $(grep -E "^@USER@:" /etc/passwd > /dev/null ) ; then
                useradd -d "/var/lib/@USER@" -M -r -c "Centreon-engine Admin" -g "@GROUP@" -s /bin/bash "@USER@"
            fi
	fi

	if [ -f "/etc/init.d/centengine" ]; then
	    update-rc.d centengine defaults
	fi

        # add permission plugins
        if [ -d "/usr/lib/nagios/plugins" ]; then
            chown @USER@:@GROUP@ /usr/lib/nagios/plugins/*
        fi

        if [ -f "/usr/lib/nagios/plugins/check_dhcp" ]; then
            chown root /usr/lib/nagios/plugins/check_dhcp
            chmod u+s /usr/lib/nagios/plugins/check_dhcp
        fi

        if [ -f "/usr/lib/nagios/plugins/check_icmp" ]; then
            chown root /usr/lib/nagios/plugins/check_icmp
            chmod u+s /usr/lib/nagios/plugins/check_icmp
        fi

        #add permission directory /etc/centreon-engine
        if [ -d "/etc/centreon-engine" ]; then
            chown -R centreon-engine:centreon-engine /etc/centreon-engine
            chmod 775 /etc/centreon-engine
        fi

        #test /var/lib/centreon-engine/rw
        if [ ! -d "/var/lib/centreon-engine/rw" ]; then
            mkdir -p /var/lib/centreon-engine/rw
        fi

        #add permission /var/lib/centreon-engine/rw
        if [ -d "/var/lib/centreon-engine/rw" ]; then
            chown -R @USER@:@GROUP@ /var/lib/centreon-engine
        fi

        #test /var/log/centreon-engine
        if [ ! -d "/var/log/centreon-engine" ]; then
            mkdir -p /var/log/centreon-engine
        fi

        #test /var/log/centreon-engine/archives
        if [ ! -d "/var/log/centreon-engine/archives" ]; then
            mkdir -p /var/log/centreon-engine/archives
        fi


        #test /var/log/centreon-engine/centengine.log
        if [ ! -f "/var/log/centreon-engine/centengine.log" ]; then
            touch /var/log/centreon-engine/centengine.log
        fi

        #add permission /var/log/centreon-engine
        if [ -d "/var/log/centreon-engine" ]; then
            chown -R @USER@:@GROUP@ /var/log/centreon-engine
        fi
	;;

    abort-upgrade|abort-remove|abort-deconfigure)
	;;

    *)
	echo "$0: called with unknown argument" 1>&2
	exit 0
	;;
esac

exit 0
